# intellij/03-intellij.conf
[program:intellij]
command=%(ENV_NVC_HOME)s/intellij.sh %(ENV_DISPLAY_1)s
startretries=%(ENV_SUPERVISOR_RETRIES)s
autostart=true