# intellij/03-intellij.conf
[program:intellij]
command=%(ENV_NVC_HOME)s/intellij.sh
autostart=true