# intellij/03-intellij.conf
[program:intellij]
command=%(ENV_HOME)s/intellij.sh
autostart=true
norestart=true