# novnc/config/supervisord/03-alacrity.conf
[program:alacritty]
command=alacritty
autorestart=false
