# novnc/config/supervisord/04-copyq.conf
[program:copyq]
command=copyq
autorestart=true
