# novnc/config/supervisord/05-terminator.conf
[program:terminator]
command=terminator
autorestart=true
