# novnc/config/supervisord/06-guaca-server.conf
[program:guacd]
command=guacd -f -l 5111 -b localhost -L debug
autorestart=true

