# novnc/config/supervisord/07-tomcat.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
autorestart=true
redirect_stderr=true