# novnc/config/supervisord/08-feh.conf
[program:feh]
command=/home/novnc/feh.sh %(ENV_ENABLE_SCREEN_1)s %(ENV_DISPLAY_1)s
autorestart=true